# docs/server/php-fpm.service
[Unit]
Description=PHP-FPM für Fiedler & Rieß Website
After=network.target

[Service]
Type=simple
ExecStartPre=-/usr/bin/podman image exists docker.io/library/php:8.3-fpm || /usr/bin/podman pull docker.io/library/php:8.3-fpm
ExecStartPre=-/usr/bin/podman rm -f php-fpm
ExecStart=/usr/bin/podman run --name php-fpm \
  --network host \
  -v /opt/fiedler-riess/www:/opt/fiedler-riess/www:ro \
  -v /opt/fiedler-riess/private:/opt/fiedler-riess/private:ro \
  -v /opt/fiedler-riess/eingaenge:/opt/fiedler-riess/eingaenge:rw \
  -v /opt/php-fpm/www.conf:/usr/local/etc/php-fpm.d/www.conf:ro \
  docker.io/library/php:8.3-fpm
ExecStop=/usr/bin/podman stop -t 10 php-fpm
ExecStopPost=-/usr/bin/podman rm -f php-fpm
Restart=always
RestartSec=5
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
